# hardware/beerrouter_openwrt_config/etc/crontabs/root

* * * * * ps | grep -q "[b]eerrouter.py" || /root/beerrouter.py &